# docker/core-mir.supervisord.conf
; Supervisor config for the core-mir merged container.
;
; Backend servers run as supervised child processes here:
;   bpm        :8004   librosa + madmom-fork
;   dsp        :8003   ruptures (change-point) + mir (MIRtoolbox features),
;                      both served by tools/python/dsp_server.py on one port
;   custom     :8005   user-uploaded Python detector sandbox
;
; Each program restarts independently on failure. supervisord stays in the
; foreground so the container's lifecycle tracks supervisor itself; the
; child stdout/stderr is teed to the container's own stdout/stderr so
; `docker logs core-mir` shows interleaved per-server output without
; needing per-service log files inside the image.

[supervisord]
nodaemon=true
; No `user=` on purpose: supervisord runs as whoever the container is
; started as. The compose `user: "${HOST_UID:-0}:${HOST_GID:-0}"` is root
; by default, but Linux users set HOST_UID to a non-root id so files written
; into data/ aren't root-owned. A hardcoded `user=root` made supervisord —
; already running as that non-root id — try to *switch to* root and abort
; with "Can't drop privilege as nonroot user". Omitting it works in both modes.
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

[program:bpm]
command=python tools/python/bpm_server.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dsp]
command=python tools/python/dsp_server.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:custom]
command=python tools/python/custom_server.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
