# docker/core-mir.supervisord.conf
; Supervisor config for the core-mir merged container.
;
; Four MIR backend servers run as supervised child processes here:
;   bpm        :8004   librosa + madmom-fork
;   mir        :8007   librosa + scipy + pyloudnorm + MIRtoolbox-parity custom
;   ruptures   :8003   librosa + ruptures change-point detection
;   custom     :8005   user-uploaded Python detector sandbox
;
; Each program restarts independently on failure. supervisord stays in the
; foreground so the container's lifecycle tracks supervisor itself; the
; child stdout/stderr is teed to the container's own stdout/stderr so
; `docker logs core-mir` shows interleaved per-server output without
; needing per-service log files inside the image.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

[program:bpm]
command=python tools/python/bpm_server.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mir]
command=python tools/python/mir_server.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ruptures]
command=python tools/python/ruptures_server.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:custom]
command=python tools/python/custom_server.py
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
